# systemd/twitter-news.timer
[Unit]
Description=Run Newsbot daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

