# systemd/twitter-news.service
[Unit]
Description=Twitter News Bot

[Service]
User=pi
Type=oneshot
ExecStart=/usr/bin/python3 /home/pi/tweetBot.py
